# systemd/aura-threat-update.timer
[Unit]
Description=Run threat intel update daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target